# btrfs-snapshot.service
[Unit]
Description=BTRFS snapshot

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/btrfs-snapshot /home /data/home-backups